# dist/linux/share/cellframe-updater.timer
[Unit]
Description=CFUPDATE

[Timer]
AccuracySec=1m
#OnBootSec=1min
OnCalendar=0/2:00:00
Persistent=true
Unit=cellframe-updater.service


[Install]
WantedBy=timers.target
