# os/macos/PKGINSTALL/postinstall
#!/bin/sh
#set -x

INSTALL_BINS="/Applications/CellframeNode.app/Contents/MacOS/"
INSTALL_RES=/Applications/CellframeNode.app/Contents/Resources
USER_RES="/Library/Application Support/CellframeNode" 

mkdir -p "${USER_RES}/"

# Copy configuration files
cp -r "${INSTALL_RES}/etc" "${USER_RES}/"
cp -r "${INSTALL_RES}/share" "${USER_RES}/"

mkdir -p "${USER_RES}/var/run/"

# Check for potential duplicate network configs before running setup
echo "[!] Checking for network config duplicates before setup"
NETWORK_DIR="${USER_RES}/etc/network"
if [ -d "${NETWORK_DIR}" ]; then
    # Use find to avoid shell glob expansion issues when no .cfg files exist
    find "${NETWORK_DIR}" -maxdepth 1 -name "*.cfg" -type f | while read -r cfg_file; do
        if [ -f "${cfg_file}" ]; then
            basename_cfg=$(basename "${cfg_file}" .cfg)
            dis_file="${cfg_file}.dis"
            if [ -f "${dis_file}" ]; then
                echo "[!] Warning: Found both ${cfg_file} and ${dis_file}, removing .dis file to prevent conflicts"
                rm -f "${dis_file}"
            fi
        fi
    done
fi

# Run configuration setup
${INSTALL_BINS}/cellframe-node-config -i "${USER_RES}/share/default.setup"

# Final check and cleanup for duplicate configs after setup
echo "[!] Post-setup cleanup: checking for network config duplicates"
if [ -d "${NETWORK_DIR}" ]; then
    # Use find to avoid shell glob expansion issues when no .cfg files exist
    find "${NETWORK_DIR}" -maxdepth 1 -name "*.cfg" -type f | while read -r cfg_file; do
        if [ -f "${cfg_file}" ]; then
            basename_cfg=$(basename "${cfg_file}" .cfg)
            dis_file="${cfg_file}.dis"
            if [ -f "${dis_file}" ]; then
                echo "[!] Warning: Found duplicate configs after setup: ${cfg_file} and ${dis_file}"
                echo "[!] Keeping user's enabled config and removing duplicate disabled config"
                rm -f "${dis_file}"
            fi
        fi
    done
fi

echo "[!] Set cfg permissions"
find "${USER_RES}/" -type d -exec chmod 777 {} \;
find "${USER_RES}/" -type f -exec chmod 666 {} \;

echo "[!] Copy daemon plist"
cp ${INSTALL_RES}/com.demlabs.cellframe-node.plist /Library/LaunchDaemons/
cp ${INSTALL_RES}/com.demlabs.cellframe-diagtool.plist /Library/LaunchDaemons/

${INSTALL_BINS}/cellframe-node-config -e service enable
launchctl load -w /Library/LaunchDaemons/com.demlabs.cellframe-diagtool.plist
