# os/debian/postinst
#!/bin/bash -e

DAP_APP_NAME="cellframe-node"
DAP_PREFIX="/opt/$DAP_APP_NAME"

# Check profile.d symlink
[ -e "/etc/profile.d/$DAP_APP_NAME.sh" ] || ln -sf $DAP_PREFIX/share/profile.d/$DAP_APP_NAME.sh /etc/profile.d/$DAP_APP_NAME.sh

#unconditianly create directories on pkg instllation
echo "[*] Creating run, lib, var dirs...."
mkdir -p $DAP_PREFIX/var/{run,lib/wallet,lib/global_db,var/plugins} || true

if [ -e /.dockerenv ]
then
    echo "[ ] No logrotate in docker-container"
else
    echo "[!] Enabling logrotate" 
    chmod 644  $DAP_PREFIX/share/logrotate/$DAP_APP_NAME || true
    chmod 644  $DAP_PREFIX/share/logrotate/logrotate.timer || true
    systemctl enable $DAP_PREFIX/share/logrotate/logrotate.timer || true
fi

#set permissions
echo "[!] Set exec permissions"
chmod +x $DAP_PREFIX/bin/$DAP_APP_NAME-cli $DAP_PREFIX/python/bin/pip3* $DAP_PREFIX/python/bin/python3* || true

echo "[!] Run cellframe-config for configuration..."
$DAP_PREFIX/bin/cellframe-node-config -i $DAP_PREFIX/share/default.setup

echo "[!] Set cfg permissions"
#set rwo permissions to configs
chmod 666 $(find ${DAP_PREFIX}/etc/ -type f)
#set rwx permissions to dirs
chmod 777 $(find ${DAP_PREFIX}/etc/ -type d)

echo "[!] Starting up cellframe-node"
systemctl enable $DAP_PREFIX/share/$DAP_APP_NAME.service || true
systemctl start cellframe-node 

if [ -e "$DAP_PREFIX/bin/cellframe-diagtool" ]; then
    echo "[!] Starting up cellframe-diagtool"
    systemctl --system enable $DAP_PREFIX/share/cellframe-diagtool.service || true
    systemctl start cellframe-diagtool
fi

if [ -e "$DAP_PREFIX/share/cellframe-updater.service" ]; then
    echo "[!] Starting up cellframe-updater"
    systemctl --system enable $DAP_PREFIX/share/cellframe-updater.service || true
    systemctl --system enable $DAP_PREFIX/share/cellframe-updater.timer || true
    systemctl start cellframe-updater
    systemctl start cellframe-updater.timer
fi


echo "[!] Done"
