# dist.linux/share/cellframe-updater.timer
[Timer]
OnCalendar=*:1/10
Persistent=true
Unit=cellframe-updater.service