# dist.linux/share/cellframe-node.service
[Unit]
Description=Cellframe Node
After=network.target

[Service]
WorkingDirectory=/opt/cellframe-node
ExecStart=/opt/cellframe-node/bin/cellframe-node &
ExecStop=/bin/kill -SIGTERM $MAINPID
Restart=always
User=cellframe-node
Group=cellframe-node
RestartSec=10

[Install]
WantedBy=multi-user.target
