# dist.linux/share/cellframe-updater.timer
[Timer]
OnCalendar=0/2:00:00
Persistent=true
unit=cellframe-updater